# online.conf
[program:m-dulo]
command=/opt/myapp/online/online
directory=/opt/myapp/online
autostart=true
autorestart=true
stderr_logfile=/var/log/online.err.log
stdout_logfile=/var/log/online.out.log
